# configs/etc/systemd/system/samba-restart.service
[Unit]
Description=Restart the samba service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart smbd.service
